# benchmark-bot.service
[Unit]
Description = isucon7q bench

[Service]
User = isucon
Group = isucon
ExecStart = /home/isucon/go/src/github.com/tukejonny/bot/bot
Restart = always
Type = simple
LimitNOFILE = 65536

[Install]
WantedBy = multi-user.target
